# docker/supervisor.sidekiq.conf
[supervisord]
nodaemon=true

[program:sidekiq]
user=root
directory=/rails_app
command=/rails_app/docker/start_sidekiq.sh
redirect_stderr=true
autostart=true
autorestart=true
