# docker/supervisor.sidekiq.conf
[supervisord]
nodaemon=true

[program:sidekiq]
user=root
directory=/rails_app
command=sidekiq
redirect_stderr=true
autostart=true
autorestart=true
