# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:seven_ten]
user=root
directory=/rails_app
command=/rails_app/docker/start.sh
redirect_stderr=true
autostart=true
autorestart=true
